# templates/default/redis_upstart.conf.erb
description "redis server"

start on runlevel [23]
stop on shutdown

setuid redis
setgid redis

expect fork
respawn

exec <%= node['redis']['source']['prefix'] %>/bin/redis-server /etc/redis/<%= node['redis']['port'] %>.conf

kill timeout 10
